# deploy/voxel-fps.service
[Unit]
Description=Voxel FPS Web App
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/tamapet
EnvironmentFile=/root/tamapet/.env
ExecStart=/bin/sh -lc '/root/tamapet/.venv/bin/uvicorn app.fps_main:app --host 0.0.0.0 --port ${FPS_PORT:-18428}'
Restart=always
RestartSec=3
TimeoutStopSec=10
KillSignal=SIGTERM
PrivateTmp=true

[Install]
WantedBy=multi-user.target
